# node_filesystem/etc/rc.local
#!/bin/sh
env
# Run Confd in onetime mode, to ensure that we have a working config in place to allow bird(s) and
# felix to start.
./confd -confdir=. -onetime -node ${ETCD_AUTHORITY}